# docker-service/config/supervisor.d/vmss-pro/td-agent.conf
[program:td-agent]
priority=150
directory=/var/log
command=/usr/sbin/td-agent
autostart=true
autorestart=true
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
